# deploy/vapor-terminal.service
# Vapor - unit del server WebSocket del terminale.
# Procedura di deploy completa e variabili d'ambiente: vedi HOWTO.md.
[Unit]
Description=Vapor - server WebSocket terminale PTY per Incus
After=network.target incus.service
Wants=incus.service

[Service]
Type=simple
# Utente con accesso al socket di Incus (gruppo incus-admin).
User=vapor
Group=incus-admin
WorkingDirectory=/opt/vapor
ExecStart=/usr/bin/php /opt/vapor/bin/terminal-server.php
Restart=on-failure
RestartSec=2

# Configurazione: sorgente unica nel file .env (vedi .env.example).
# Lo stesso file è letto anche dalla CLI e dal web server.
EnvironmentFile=/opt/vapor/.env

# Hardening di base.
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true

[Install]
WantedBy=multi-user.target
